# mirosim_servers/mirosim.service
[Unit]
Description=Mirosim Service

[Service]
User=pisrv
ExecStart=/home/pisrv/docker_home/startSims.sh

[Install]
WantedBy=multi-user.target
